# docker/supervisord.conf
# supervisord runs as PID 1 (non-root, appuser) and manages nginx + N
# single-worker uvicorns. Logs go to the container stdout/stderr (prefixed per
# program). On `docker stop` (SIGTERM) it stops the uvicorns first (priority 20)
# before nginx (priority 10): each uvicorn receives SIGTERM → runs its lifespan
# shutdown → final leaderboard flush of pending counts.
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
# run as the container's non-root user (set by Dockerfile USER); supervisor
# inherits it and spawns nginx + uvicorn as the same user.

# Local control socket so `docker exec <c> supervisorctl status` works for ops.
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:nginx]
command=/usr/sbin/nginx -c /app/docker/nginx.conf -g "daemon off;"
autostart=true
autorestart=true
startsecs=1
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:uvicorn-1]
command=/app/.venv/bin/ai-speech-shadowing serve --host 127.0.0.1 --port 8001
directory=/app
environment=WORKERS="1"
autostart=true
autorestart=true
startsecs=3
priority=20
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:uvicorn-2]
command=/app/.venv/bin/ai-speech-shadowing serve --host 127.0.0.1 --port 8002
directory=/app
environment=WORKERS="1"
autostart=true
autorestart=true
startsecs=3
priority=20
stopsignal=TERM
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
redirect_stderr=true
